# scripts/systemd/nexus-scheduled-reports.timer
[Unit]
Description=Nexus rapports planifiés — toutes les 15 min
Requires=nexus-scheduled-reports.service

[Timer]
OnCalendar=*:0/15
RandomizedDelaySec=30
Persistent=true
Unit=nexus-scheduled-reports.service

[Install]
WantedBy=timers.target
