# scripts/systemd/nexus-weekly-digest.service
[Unit]
Description=Nexus weekly digest — résumé de semaine vendredi 17h
After=network.target nexus.service

[Service]
Type=oneshot
User=cetix
WorkingDirectory=/opt/nexus
EnvironmentFile=/opt/nexus/.env
ExecStart=/usr/bin/env npx tsx src/workers/weekly-digest-worker.ts
StandardOutput=append:/var/log/nexus-weekly-digest.log
StandardError=append:/var/log/nexus-weekly-digest.log
TimeoutStartSec=300
